# templates/etc/profile.d/spark.sh
export PATH=$PATH:{{ spark_usr_dir }}/bin/
